# packaging/nfpm/pivirt-agent.prerm
#!/bin/bash
set -e

# Stop the service
systemctl stop pivirt-agent.service

# Disable the service
systemctl disable pivirt-agent.service

exit 0
